# config/mate.service
[Unit]
Description=Mate Local Assistant Daemon
After=network.target ollama.service
Requires=ollama.service

[Service]
Type=simple
ExecStart=/usr/local/bin/mate-daemon
User=mate
Group=mate
Restart=always
RestartSec=5
KillMode=process    # <<< important: only kill main process, not the whole cgroup
TimeoutStopSec=20   # give daemon time to shutdown gracefully

[Install]
WantedBy=multi-user.target
